# board/hermes-pi5/rootfs-overlay/etc/systemd/system/gpsd.socket.d/no-ipv6.conf
[Socket]
ListenStream=
ListenStream=/run/gpsd.sock
ListenStream=127.0.0.1:2947
