# package/hermes-music/hermes-music-bootstrap.service
[Unit]
Description=Hermes Music — first-boot bootstrap (venv + pip install)
After=network-online.target
Wants=network-online.target
Before=hermes-music.service
ConditionPathExists=!/var/lib/hermes-music/.bootstrap-done

[Service]
Type=oneshot
ExecStart=/usr/lib/hermes-music/bootstrap.sh
RemainAfterExit=yes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
