# board/hermes-pi5/rootfs-overlay/etc/systemd/system/wpa_supplicant.service.d/hermes.conf
[Service]
# Hermes override: load persistent config (update_config=1) so networks
# added via D-Bus AddNetwork survive reboots, and pre-create the wlan0
# interface so elise's NetworkController finds it without CreateInterface
# (which would not pass ConfigFile).
ExecStart=
ExecStart=/usr/sbin/wpa_supplicant -u -c /etc/wpa_supplicant/wpa_supplicant.conf -i wlan0
