# board/hermes-pi5/rootfs-overlay/etc/systemd/system/wireplumber.service.d/state-dir.conf
[Service]
# WirePlumber persists per-device routing decisions (default sink,
# preferred BT codec) under $XDG_STATE_HOME. The packaged user
# `pipewire` has HOME=/, so wp falls back to /.local/state/wireplumber
# which it can't write — every boot loses pairing/routing state.
# Pin both vars to a real, writable per-service state dir.
StateDirectory=wireplumber
Environment=XDG_STATE_HOME=/var/lib/wireplumber
Environment=HOME=/var/lib/wireplumber
